# crontab
*/5 * * * * python3 /script/ddns_bot.py
